# audio-topology-profile/audio-topology.service
[Unit]
Description=Audio topology: IRQ affinity tuning
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/audio-topology-profile/irq-affinity.sh

[Install]
WantedBy=multi-user.target
